# supervisor/websocket.conf
[program:news_websocket]
command=/usr/bin/python3 /var/www/ai-news-automation/src/update_server.py
directory=/var/www/ai-news-automation
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/supervisor/news_websocket.err.log
stdout_logfile=/var/log/supervisor/news_websocket.out.log
environment=
    NEWS_WS_HOST="newssurgeai.com",
    NEWS_WS_PORT="443",
    NEWS_API_KEY="%(ENV_NEWS_API_KEY)s",
    NODE_ENV="production"

[supervisord]
nodaemon=true

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
